# scripts/miloco-agent-supervisor.conf.example
#!/usr/bin/env bash
# Fork 专属：supervisor 示例（7×24 运行 miloco-agent Sidecar）
# 用法：复制到 /etc/supervisor/conf.d/ 或 include 到现有 supervisor 配置后 supervisorctl reread && update

[program:miloco-agent]
command=bash %(ENV_MILOCO_REPO)s/scripts/miloco-agent-run.sh
directory=%(ENV_MILOCO_REPO)s
user=%(ENV_USER)s
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=20
stdout_logfile=%(ENV_MILOCO_HOME)s/logs/miloco-agent.log
stderr_logfile=%(ENV_MILOCO_HOME)s/logs/miloco-agent.err.log
environment=MILOCO_HOME="%(ENV_MILOCO_HOME)s",MILOCO_AGENT_VENV="%(ENV_MILOCO_REPO)s/miloco-agent/.venv"
